# generic-python-flask-supervisord.conf
[supervisord]
nodaemon=true

[program:generic-python-flask]
command=/usr/local/bin/python /application/__init__.py
priority=900
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
username=www-data
autorestart=true
